# live os/config/includes.chroot/etc/systemd/system/ai-live-password.service
[Unit]
Description=Generate per-boot password for live SSH user
After=live-config.service
Wants=live-config.service
Before=ssh.service getty@tty1.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ai-generate-live-password
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
