# live os/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Unit]
Wants=live-config.service
After=live-config.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin live --noclear %I $TERM
